# scripts/RunMailToLeahXbpCsv.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2025-12-01T08:00:00</Date>
    <Author>EC2-Window-FIPS-enabled\Administrator</Author>
    <Description>Run daily report job at 4:01 AM</Description>
  </RegistrationInfo>

  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2025-12-01T23:40:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <StartWhenAvailable>true</StartWhenAvailable>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
  </Settings>

  <!-- <Actions Context="Author">
    <Exec>
      <Command>C:\Users\Administrator\Documents\ees\envEes\Scripts\python.exe</Command>
      <Arguments>manage.py run_mail_to_leah_xbp_csv</Arguments>
      <WorkingDirectory>C:\Users\Administrator\Documents\ees</WorkingDirectory>
    </Exec>
  </Actions> -->
  <Actions Context="Author">
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c "set PYTHONPATH=C:\Users\Administrator\Documents\ees &amp;&amp; cd /d C:\Users\Administrator\Documents\ees &amp;&amp; C:\Users\Administrator\Documents\envEes\Scripts\python.exe manage.py run_mail_to_leah_xbp_csv"</Arguments>
      <WorkingDirectory>C:\Users\Administrator\Documents\ees</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
