# scripts/RunMidNightCleanUp.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>EC2-Window-FIPS-enabled\Administrator</Author>
    <Description>Midnight cleanup (adj + sup removal) at 23:50</Description>
  </RegistrationInfo>

  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2025-12-01T23:50:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <StartWhenAvailable>true</StartWhenAvailable>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
  </Settings>

  <!-- <Actions Context="Author">
    <Exec>
      <Command>C:\Users\Administrator\Documents\ees\envEes\Scripts\python.exe</Command>
      <Arguments>manage.py run_midnight_cleanup</Arguments>
      <WorkingDirectory>C:\Users\Administrator\Documents\ees</WorkingDirectory>
    </Exec>
  </Actions> -->
  <Actions Context="Author">
  <Exec>
    <Command>cmd.exe</Command>
    <Arguments>/c "set PYTHONPATH=C:\Users\Administrator\Documents\ees &amp;&amp; cd /d C:\Users\Administrator\Documents\ees &amp;&amp; C:\Users\Administrator\Documents\envEes\Scripts\python.exe manage.py run_midnight_cleanup"</Arguments>
    <WorkingDirectory>C:\Users\Administrator\Documents\ees</WorkingDirectory>
  </Exec>
</Actions>
</Task>
